# home/private_dot_config/systemd/user/rclone-photo.timer
[Unit]
Description=Run Photo sync every day

[Timer]
Unit=rclone-photo.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
